# root/etc/services.d/encfs/run
#!/usr/bin/with-contenv sh
if [ $ENCFS_PWD ]; then
   /usr/bin/encfs $ENCFS_OPTS --extpass="echo $ENCFS_PWD" /src /dest
else
   /usr/bin/encfs $ENCFS_OPTS /src /dest
fi
tail -f /dev/null
