# root/etc/services.d/encfs/run
#!/usr/bin/with-contenv sh

cat "/config/passwd" | /usr/bin/encfs -S -o allow_other,nonempty /encrypted /decrypted 
tail -f /dev/null
