# root/etc/services.d/encfs/run
#!/usr/bin/with-contenv sh

cat "/config/passwd" | /usr/bin/encfs ${ENCFS_OPTS} /src /dest
tail -f /dev/null
